# conf/etc/periodic/15min/01-wordpress
#!/bin/sh


run_wp(){
  
  eval "su - nginx -s /bin/ash -c 'cd /var/www/wordpress/ && "$@"'" 
  return $? 
}

run_wp wp cron